# systemd/b2_backup.timer
[Unit]
Description=B2 backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
